# vps/second-brain-hub/deploy/com.rbedu.mrluc-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rbedu.mrluc-export</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_PATH/vps/second-brain-hub/scripts/export_vault_to_vps.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>6</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>6</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>6</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/mrluc-export.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mrluc-export.err</string>
</dict>
</plist>
